# diy-q-business-chatbot/AMI/files/ssl-cert.service
[Unit]
Description=renew ssl certificate
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ssl-cert.sh