# diy-q-business-chatbot/AMI/files/ssl-cert.timer
[Unit]
Description=runs script that renews ssl certificate

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target